# config/systemd/certbot.timer
[Unit]
Description=Run certbot daily and on boot

[Timer]
OnBootSec=1m
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target